# deploy/gale.service
[Unit]
Description=Gale static web server
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/gale
WorkingDirectory=/var/www
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
